# i3wm/.config/systemd/user/redshiftconf.service
[Unit]
Description=redshiftconf

[Service]
ExecStart=/usr/bin/redshiftconf

[Install]
WantedBy=default.target
